# fileengine.service
[Unit]
Description=FileEngine Core Service
After=network.target postgresql.service
Wants=postgresql.service

[Service]
# Type=notify: the server sends sd_notify("READY=1") once both the gRPC
# listener and the REST listener are bound. systemd will hold ordering
# (After=/Before= for downstream units) until that happens.
Type=notify
NotifyAccess=main
User=fileengine
Group=fileengine
ExecStart=/usr/bin/fileengine_server
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=3
# Generous startup timeout — initial DB-pool warm-up + S3 init can be slow.
TimeoutStartSec=60
StandardOutput=journal
StandardError=journal
SyslogIdentifier=fileengine

# Create necessary directories
StateDirectory=fileengine
LogsDirectory=fileengine

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/fileengine /var/log/fileengine /tmp/fileengine

# Hardening settings
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
SecureBits=keep-caps
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictRealtime=true
RestrictNamespaces=true
RestrictSUIDSGID=true
RemoveIPC=true

[Install]
WantedBy=multi-user.target