# fileengine.service
[Unit]
Description=FileEngine Core Service
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=fileengine
Group=fileengine
ExecStart=/usr/bin/fileengine_server
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal
SyslogIdentifier=fileengine

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/fileengine /var/log/fileengine /tmp/fileengine

# Hardening settings
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
SecureBits=keep-caps
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictRealtime=true
RestrictNamespaces=true
RestrictSUIDSGID=true
RemoveIPC=true

[Install]
WantedBy=multi-user.target