# debian/postrm
#!/bin/bash
set -e

case "$1" in
    remove)
        # Disable the service
        systemctl disable fileengine.service >/dev/null || true
    ;;
    purge)
        # Remove the fileengine user and group
        deluser --system fileengine 2>/dev/null || true
        delgroup --system fileengine 2>/dev/null || true
    ;;
    upgrade|failed-upgrade|abort-install|abort-upgrade|disappear)
    ;;
esac

# Reload systemd daemon
if [ -d /run/systemd/system ]; then
    systemctl daemon-reload || true
fi

#DEBHELPER#

exit 0